# SOURCES/ihttpd.path
# Path unit starting ihttpd.service on password request
[Unit]
Description=The Apache IHTTP Server that will unlock crypttab device
DefaultDependencies=no
Conflicts=shutdown.target multi-user.target graphical.target
After=plymouth-start.service
Before=paths.target shutdown.target
ConditionPathExists=!/run/plymouth/pid

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
